# journal.cron
crontab -e

30 9 * * * /usr/bin/touch "$HOME/Documents/HKP_"$(date "+%Y_%m_%d")".md"
